# service/macos/com.codemonitor.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codemonitor.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/monitor-server</string>
        <string>--config</string>
        <string>/usr/local/etc/code-monitor/config.toml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/usr/local/etc/code-monitor</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/code-monitor/server.out.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/code-monitor/server.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
